# addons_sdk/SQUID/files/etc/cron.d/squid_rotate

# rotate squid logs
00 0 * * *	root pidof squid &> /dev/null && /opt/squid/sbin/squid -k rotate
